# systemd/onecmd.service
[Unit]
Description=OneCmd Telegram terminal manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@ROOT_DIR@
ExecStart=@ROOT_DIR@/run.sh
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=20

# Hardening (safe defaults for user service)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=@ROOT_DIR@ %h/.onecmd

[Install]
WantedBy=default.target
